# schedule.cron
0 * * * * python3 /home/opc/AutoScale/AutoScaleALL.py
