# schedule.cron
0 * * * * python3 /home/opc/OCI-AutoScale/AutoScaleALL.py
